# ansible/roles/es-agent/templates/elastic-agent.service.j2
[Unit]
Description=Elastic Agent
After=network.target

[Service]
ExecStart=/usr/share/elastic-agent/elastic-agent run
Restart=always
User=root

[Install]
WantedBy=multi-user.target